# .docker/supervisor/mysql.conf
[program:mysqld]
command=/opt/docker/server/start-mysql.sh
numprocs=1
autostart=true
autorestart=true